# mcp-gateway/deploy/mcp-gateway.service
[Unit]
Description=Huawei Cloud MCP Gateway (multi-server SSE on one port)
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=mcp-gateway
Group=mcp-gateway
WorkingDirectory=/opt/mcp-servers
# Single .env file at the workspace root — contains all shared credentials
# (AK/SK, project ID, JWT key) plus per-service region/log config.
EnvironmentFile=/etc/mcp-gateway/.env
# Use the workspace-root start.sh (loads .env + delegates to mcp-gateway CLI).
# manifest.yaml lives at the workspace root alongside .env and start.sh.
ExecStart=/opt/mcp-servers/start.sh \
    --manifest /opt/mcp-servers/manifest.yaml \
    --host ${MCP_GATEWAY_HOST} \
    --port ${MCP_GATEWAY_PORT} \
    --log-level ${MCP_GATEWAY_LOG_LEVEL}
Restart=on-failure
RestartSec=5s
# Conservative resource caps. Tune to your environment.
LimitNOFILE=65536
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log/mcp-gateway

[Install]
WantedBy=multi-user.target
